# etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user start-up script.
#

# Set radeon card into low power mode
echo mid > /sys/class/drm/card0/device/power_profile

echo 256 > /sys/block/sda/queue/read_ahead_kb

# change the scheduler queue length
echo 1024 > /sys/block/sda/queue/nr_requests

# Set to never or madvise to avoid lag during heavy IO
# http://lwn.net/Articles/467328/
echo never > /sys/kernel/mm/transparent_hugepage/defrag

# Turn on audio card power save
echo Y > /sys/module/snd_hda_intel/parameters/power_save_controller
echo 1 > /sys/module/snd_hda_intel/parameters/power_save

# Disable nmi watchdog
echo 0 > /proc/sys/kernel/nmi_watchdog

# do not park the HDD head please
hdparm -B 254 /dev/sda
