# etc/systemd/system/cleanup.timer
[Unit]
Description=Regular pacman cache cleanup

[Timer]
OnCalendar=weekly

[Install]
WantedBy=multi-user.target
